# bonny-server/templates/{project}.service
[Unit]
Description={project}

[Service]
ExecStart=racket /home/racket/{project}/server.rkt 8000
WorkingDirectory=/home/racket/{project}
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier={project}
User=racket
Group=racket
Environment=
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
